# its-a-plane-python/setup/plane-tracker.service
[Unit]
Description=Plane Tracker RGB LED Matrix Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Secrets loaded from a root-only environment file
EnvironmentFile=/etc/plane-tracker.env
WorkingDirectory=/home/robk/plane-tracker-rgb-pi/its-a-plane-python

# Fix permissions on writable data files/dirs before starting
ExecStartPre=/bin/bash -c 'chmod -R 777 /home/robk/plane-tracker-rgb-pi/its-a-plane-python/web/static/maps 2>/dev/null; touch /home/robk/plane-tracker-rgb-pi/its-a-plane-python/close.txt /home/robk/plane-tracker-rgb-pi/its-a-plane-python/farthest.txt && chmod 666 /home/robk/plane-tracker-rgb-pi/its-a-plane-python/close.txt /home/robk/plane-tracker-rgb-pi/its-a-plane-python/farthest.txt'

ExecStart=/usr/bin/python3 /home/robk/plane-tracker-rgb-pi/its-a-plane-python/its-a-plane.py
Restart=on-failure
RestartSec=10
# Run as root for GPIO / LED matrix access
User=root
Group=root

[Install]
WantedBy=multi-user.target
