# its-a-plane-python/setup/plane-tracker.service
[Unit]
Description=Plane Tracker RGB LED Matrix Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Secrets and config loaded from environment file
EnvironmentFile=/etc/plane-tracker.env
WorkingDirectory=__REPO_DIR__/its-a-plane-python
Environment=PLANE_TRACKER_DATA_DIR=/var/lib/plane-tracker

# Ensure writable data directory exists (owned by root since drop_privileges=False)
ExecStartPre=/bin/bash -c "mkdir -p /var/lib/plane-tracker/maps && chown -R root:root /var/lib/plane-tracker"

ExecStart=__REPO_DIR__/.venv/bin/python3 __REPO_DIR__/its-a-plane-python/its-a-plane.py
Restart=on-failure
RestartSec=10
# Run as root for GPIO / LED matrix access
User=root
Group=root

[Install]
WantedBy=multi-user.target
